# crontab
# 
# m h  dom mon dow   command
* 22 * * * /root/domainscan.sh >> /var/log/domainscan.log
